# supervisord.conf
# Security model:
# - supervisord runs as root (required to manage nginx binding to port 80)
# - nginx master runs as root, worker processes drop to nobody (default nginx behavior)
# - uvicorn (backend) drops to appuser via user= directive
# - HEALTHCHECK uses curl against nginx on localhost:80

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:backend]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app/backend
user=appuser
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
